# data/com.github.darkoverlordofdata.catlock.cron
0 */6 * * * com.github.darkoverlordofdata.badabing --update >/dev/null 2>&1